# config/systemd/llama-primary.service
[Unit]
Description=ik_llama.cpp Server - Primary Model
After=network.target
Documentation=https://github.com/ikawrakow/ik_llama.cpp

[Service]
Type=simple
User=root
# __MODEL_PATH__ is replaced by install.sh with the actual best model path
ExecStart=/usr/local/bin/ik-llama-server \
  --model __MODEL_PATH__ \
  --threads __THREADS__ \
  --numa distribute \
  --mlock \
  --ctx-size 8192 \
  --ctk q8_0 --ctv q4_0 \
  --batch-size 512 \
  --ubatch-size 256 \
  --host 0.0.0.0 \
  --port 8080 \
  --metrics
Restart=always
RestartSec=10
LimitMEMLOCK=infinity
# Protect system
ProtectSystem=strict
ReadWritePaths=/opt/models /tmp
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
